# ui/web/com.openhealth.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Автозапуск дашборда OpenHealth как фонового сервиса (не падает, всегда на localhost:8770).
     Установка: см. DESKTOP.md. __WEB_DIR__ заменить на абсолютный путь до ui/web. -->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.openhealth.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>http.server</string>
    <string>8770</string>
  </array>
  <key>WorkingDirectory</key><string>__WEB_DIR__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/tmp/openhealth-dashboard.log</string>
  <key>StandardErrorPath</key><string>/tmp/openhealth-dashboard.log</string>
</dict>
</plist>
